# deploy/roles/mumble-bot/templates/mumble-bot.service.j2
[Unit]
Description=Mumble Bot
After=network.target

[Service]
TimeoutStartSec=5
#User=prior
WorkingDirectory=/home/prior/applications/mumble-bot/
ExecStart=/home/prior/applications/mumble-bot/node-v{{ node_version }}-linux-x64/bin/node server serve --config-file config.yml
TimeoutStopSec=10
Restart=always
RestartSec=10

[Install]
WantedBy=multi.target
